# docs/ops/rookery.service
# Template consumed by `rookery init`, which writes the rendered unit to
# ./rookery.service in the working directory with User= filled in from
# --user / `whoami`. `sudo rookery install` then copies the rendered file
# to /etc/systemd/system/rookery.service and runs daemon-reload. See
# docs/adr/ADR-0033.md.

[Unit]
Description=rookery mail server
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/rookery
ExecStart=/opt/rookery/rookery start --prod -d
ExecStop=/opt/rookery/rookery stop
User=YOUR_UNIX_USER

[Install]
WantedBy=multi-user.target
